# scripts/deploy/tontine-api.service
[Unit]
Description=Tontine API Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/apps/tontine
ExecStart=/root/.volta/bin/node dist/main.js
Restart=always
RestartSec=10
StandardOutput=append:/root/apps/tontine/logs/app.log
StandardError=append:/root/apps/tontine/logs/app.log
Environment=NODE_ENV=production
Environment=SENTRY_DSN=%i
Environment=SENTRY_ORG=%i
Environment=SENTRY_PROJECT=%i

# Gestion des signaux
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Sécurité
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/root/apps/tontine/logs

[Install]
WantedBy=multi-user.target
